# build/services/hugo/finish
#!/bin/bash -l
# WARNING: This file has been auto-generated. DO NOT EDIT: WARNING
# Created on 2020-02-18T21:11:12+1100, using template:finish.tmpl and json:gearbox.json

# Required to avoid services starting too quickly.
sleep 5

c_ok "Stopped."

