# build/services/hugo/run
#!/bin/bash

echo "Hugo running"
export HUGO_BIN="/usr/local/bin/hugo-extended"
export HUGO_SRCDIR="/home/gearbox/projects/default/hugo"
export HUGO_PUBLISHDIR="/home/gearbox/projects/default/www"

if [ ! -d "${HUGO_SRCDIR}" ]
then
	rsync -HvaxP /hugo/ ${HUGO_SRCDIR}/
fi

cd /home/gearbox/projects/default/hugo
if [ -f config.toml ]
then
	echo "# Config file found."

elif [[ config.toml.tmpl -nt config.toml ]]
then
	echo "# Template hasn't changed - no update."

else
	echo "# Template has changed - updating."
	# JsonToConfig -json config.toml.json -template config.toml.tmpl -out config.toml
fi

su - gearbox -s /bin/bash --preserve-environment -c "${HUGO_BIN} server -D --verbose --minify --bind=\"0.0.0.0\" --port=1313 --destination=\"${HUGO_PUBLISHDIR}\" --renderToDisk"

